# automation/com.ecomhq.bosta_daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ecomhq.bosta_daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/seifeltohamy/Desktop/Zen/sku-app/.venv/bin/python</string>
    <string>/Users/seifeltohamy/Desktop/Zen/sku-app/automation/bosta_daily.py</string>
  </array>

  <!-- Run every day at 07:00 AM -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>   <integer>7</integer>
    <key>Minute</key> <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/bosta_daily_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/bosta_daily_stderr.log</string>

  <!-- Do NOT run immediately on load, only on schedule -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
